# supervisord.conf
[supervisord]
nodaemon=true

[program:g6]
command=/usr/bin/php -q /app/run-g6.php

[program:sshd]
command=/usr/sbin/sshd -D
